# overlays/System/Library/LaunchDaemons/io.github.gershwin-desktop.loginwindow.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Gershwin login greeter — replaces the console getty as the boot login UI.
  /System/Library/Scripts/LoginWindow.sh sources GNUstep.sh, sets DISPLAY=:0,
  and execs /System/Library/CoreServices/Applications/LoginWindow.app/LoginWindow.
  Session state lives in /Local/Library/Preferences/LoginWindow.plist.
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.github.gershwin-desktop.loginwindow</string>
    <key>ProgramArguments</key>
    <array>
        <string>/System/Library/Scripts/LoginWindow.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/var/log/LoginWindow.log</string>
    <key>StandardErrorPath</key>
    <string>/var/log/LoginWindow.log</string>
</dict>
</plist>
